# reliance net card/airtel-start.rules
ACTION=="add",SUBSYSTEM=="usb",ENV{DEVTYPE}=="usb_device",ATTRS{idVendor}=="12d1",ATTRS{idProduct}=="1446",RUN+="/usr/local/startAirtel.sh"
