# reliance net card/reliance-start.rules
ACTION=="add",SUBSYSTEM=="usb",ENV{DEVTYPE}=="usb_device",ATTRS{idVendor}=="19d2",ATTRS{idProduct}=="ffe9",RUN+="/usr/local/startRe.sh"
